# ctv-bot.service
[Unit]
Description=CTV Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ctv-bot
ExecStart=/home/ubuntu/ctv-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=DATA_DIR=/home/ubuntu/ctv-bot-data

[Install]
WantedBy=multi-user.target
